# supervisor-app.ini
[program:app-uwsgi]
command = uwsgi --ini /srv/app/uwsgi.ini --logger file:/var/log/uwsgi.log

[program:nginx-app]
command = nginx
           
